# configs/systemd/offlineimap@.service
[Unit]
Description=OfflineIMAP
After=network.target

[Service]
ExecStart=/usr/bin/offlineimap -c ~/.offlineimap/%i.rc
