# home/.config/systemd/user/offlineimap@.service
[Unit]
Description=OfflineIMAP
After=network.target

[Service]
ExecStart={{HOME}}/.local/bin/offlineimap -c ~/.offlineimap/%i.rc
