# home/.config/systemd/user/offlineimap@.timer
[Unit]
Description=OfflineIMAP Timer

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
