# home/.config/systemd/user/rclone.timer
[Unit]
Description=rclone timer

[Timer]
OnCalendar=*:0/30
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
