# home/.config/systemd/user/rclone.service
[Unit]
Description=rclone bisync

[Service]
Type=oneshot
ExecStart=/home/me/.local/bin/rclone-bisync
