# home/.config/systemd/user/reposync@.service
[Unit]
Description=RepoSync
After=network.target

[Service]
ExecStart={{BIN}}/reposync {{REPO}}/%i
